# deploy/stock-assistant-api.service
[Unit]
Description=Stock Assistant FastAPI backend
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=stockassistant
Group=stockassistant
UMask=0077
WorkingDirectory=/opt/stock-assistant/backend
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/opt/stock-assistant/backend"
Environment="SERVICE_NAME=stock-assistant-api"
Environment="LOG_FORMAT=json"
EnvironmentFile=/etc/stock-assistant/stock-assistant.env
Environment="AUTH_REQUIRED=true"
Environment="AUTH_TRUST_PROXY=true"
Environment="TASK_QUEUE_MODE=celery"
Environment="AGENT_WORKER_ENABLED=0"
Environment="DECISION_CHECK_WORKER_ENABLED=0"
Environment="AGENT_MAX_BATCH_SIZE=6"
Environment="AGENT_MAX_PENDING_RUNS=20"
Environment="REQUIRE_OBJECT_STORAGE=true"
ExecStart=/opt/stock-assistant/venv/bin/python -m uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
TimeoutStopSec=30
MemoryMax=1G
TasksMax=256
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
ReadWritePaths=/var/lib/stock-assistant

[Install]
WantedBy=multi-user.target
