# deploy/stock-assistant-backup.timer
[Unit]
Description=Daily Stock Assistant PostgreSQL backup

[Timer]
OnCalendar=*-*-* 03:25:00
RandomizedDelaySec=20m
Persistent=true
Unit=stock-assistant-backup.service

[Install]
WantedBy=timers.target

