# deploy/stock-assistant-agent-worker.service
[Unit]
Description=Stock Assistant Agent orchestration worker
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=stockassistant
Group=stockassistant
UMask=0077
WorkingDirectory=/opt/stock-assistant/backend
EnvironmentFile=/etc/stock-assistant/stock-assistant.env
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/opt/stock-assistant/backend"
Environment="SERVICE_NAME=stock-assistant-agent-worker"
Environment="LOG_FORMAT=json"
Environment="TASK_QUEUE_MODE=celery"
Environment="MALLOC_ARENA_MAX=2"
ExecStart=/opt/stock-assistant/venv/bin/celery -A task_queue:celery_app worker --queues=agent --hostname=agent@%H --concurrency=2 --pool=prefork --max-tasks-per-child=20 --without-gossip --without-mingle
Restart=always
RestartSec=5
TimeoutStopSec=90
KillSignal=SIGTERM
MemoryMax=1100M
TasksMax=256
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/stock-assistant

[Install]
WantedBy=multi-user.target
