# deploy/stock-assistant-healthcheck.service
[Unit]
Description=Stock Assistant dependency and worker health check
After=stock-assistant-api@8001.service stock-assistant-api@8002.service

[Service]
Type=oneshot
User=stockassistant
Group=stockassistant
Environment="API_REPLICA_PORTS=8001 8002"
Environment="API_REPLICA_REQUIRED=2"
ExecStart=/opt/stock-assistant/deploy/scripts/check-runtime.sh
TimeoutStartSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
