# deploy/stock-assistant-backup-verify.timer
[Unit]
Description=Weekly Stock Assistant PostgreSQL restore drill

[Timer]
OnCalendar=Sun *-*-* 04:30:00
RandomizedDelaySec=30m
Persistent=true
Unit=stock-assistant-backup-verify.service

[Install]
WantedBy=timers.target

