# deploy/stock-assistant-healthcheck.timer
[Unit]
Description=Run Stock Assistant health check every two minutes

[Timer]
OnBootSec=3m
OnUnitActiveSec=2m
AccuracySec=15s
Persistent=true
Unit=stock-assistant-healthcheck.service

[Install]
WantedBy=timers.target
